# deploy/systemd/lego-renew.service
[Unit]
Description=Renew DocxAutoFill Let's Encrypt short-lived certificate
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/docxautofill/scripts/renew-cert.sh
